# dockerfiles/webvm-underlay.openrc
#!/sbin/openrc-run

description="WebVM Ethernet-only FIPS underlay"

depend() {
    need localmount
    after modules
    before webvm-hashtree webvm-nvpn
}

start() {
    ebegin "Preparing WebVM FIPS Ethernet underlay"
    /usr/local/sbin/webvm-network && /usr/local/sbin/webvm-tun
    eend $?
}

stop() {
    return 0
}
